# ops/systemd/trader-us-discord-status.timer
[Unit]
Description=Run Trader-US Discord Status Notifier Every 6 Hours (KST)

[Timer]
# Every 6 hours based on KST. (Cooldown is enforced by the notifier itself.)
OnCalendar=*-*-* 03:00:00 Asia/Seoul
OnCalendar=*-*-* 09:00:00 Asia/Seoul
OnCalendar=*-*-* 15:00:00 Asia/Seoul
OnCalendar=*-*-* 21:00:00 Asia/Seoul
Persistent=true

[Install]
WantedBy=timers.target
